# agent/debian/python-lmc-base.postinst
#!/bin/sh

set -e

case "$1" in
    configure)
	[ -x /etc/init.d/lmc-agent ] && invoke-rc.d lmc-agent start
	;;

    abort-upgrade|abort-remove|abort-deconfigure)

	;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
	;;
esac

#DEBHELPER#

exit 0
