# agent/debian/python-lmc-mail.prerm
#!/bin/sh

# prerm script

set -e

case "$1" in
    remove)
        [ -x /etc/init.d/lmc-agent ] && invoke-rc.d lmc-agent stop
        ;;
    upgrade|deconfigure)
	;;
    failed-upgrade)
        ;;
    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
