# debian/pulse2-scheduler.postinst
#!/bin/sh
set -e

if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d pulse2-scheduler stop || true
else
        /etc/init.d/pulse2-scheduler stop || true
fi

if [ ! -d  "/var/lib/pulse2/imaging/inventories" ]; then
	mkdir -p /var/lib/pulse2/imaging/inventories
fi

chmod -Rv 777 /var/lib/pulse2/imaging/inventories

#DEBHELPER#

exit 0
