# services/clients/agent/linux/pulse2-agent.service
[Unit]
Description=Pulse2 Agent
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/pulse2-agent

[Install]
WantedBy=multi-user.target
